# docker/prod/worker-supervisor.conf
[program:push_notification_sender_android]
command=php app/console rabbitmq:consumer -l 150 -w -m 0 push_notification_sender_android -vv --no-ansi
directory=/www/awardwallet
umask=002
startsecs=30
startretries=100
autorestart=true
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
process_name=%(program_name)s-%(process_num)s
numprocs=2

[program:push_notification_sender_ios]
command=php app/console rabbitmq:consumer -l 150 -w -m 0 push_notification_sender_ios -vv --no-ansi
directory=/www/awardwallet
umask=002
startsecs=30
startretries=100
autorestart=true
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
process_name=%(program_name)s-%(process_num)s
numprocs=2

[program:push_notification_sender_pushy]
command=php app/console rabbitmq:consumer -l 150 -w -m 0 push_notification_sender_pushy -vv --no-ansi
directory=/www/awardwallet
umask=002
startsecs=30
startretries=100
autorestart=true
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
process_name=%(program_name)s-%(process_num)s
numprocs=1

[program:push_notification_device]
command=php app/console rabbitmq:consumer -l 150 -w -m 0 push_notification_device -vv --no-ansi
directory=/www/awardwallet
umask=002
startsecs=30
startretries=30
autorestart=true
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
process_name=%(program_name)s-%(process_num)s
numprocs=1

[program:trip_alerts_updater]
command=php app/console rabbitmq:consumer -l 500 -w -m 0 trip_alerts_updater -vv --no-ansi
directory=/www/awardwallet
umask=002
startsecs=30
startretries=30
autorestart=true
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
process_name=%(program_name)s-%(process_num)s
numprocs=1

[program:itinerary_notification]
command=php app/console rabbitmq:consumer -l 128 -w -m 0 itinerary_notification -vv --no-ansi
directory=/www/awardwallet
umask=002
startsecs=30
startretries=3000
autorestart=true
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
process_name=%(program_name)s-%(process_num)s
numprocs=1

[program:flight_info_updater]
command=php app/console rabbitmq:consumer -l 100 -m 0 -w flight_info_updater -vv --no-ansi
directory=/www/awardwallet
umask=002
startsecs=30
startretries=3000
autorestart=true
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
process_name=%(program_name)s-%(process_num)s
numprocs=1

[program:loyalty_callback_processor]
command=php app/console rabbitmq:consumer -l 150 -m 2000 -w loyalty_callback_processor -vv --no-ansi
directory=/www/awardwallet
umask=002
startsecs=30
startretries=3000
autorestart=true
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
process_name=%(program_name)s-%(process_num)s
numprocs=2

# to process old messages, scheduled pushes, delete this job after (1?) month
[program:async_processor]
command=php app/console rabbitmq:consumer -w async_processor -l 150 -vv --no-ansi
directory=/www/awardwallet
umask=002
startsecs=30
startretries=30000
autorestart=true
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
process_name=%(program_name)s-%(process_num)s
numprocs=1

[program:async_processor_2]
command=php app/console rabbitmq:consumer -w async_processor_2 -l 150 -vv --no-ansi
directory=/www/awardwallet
umask=002
startsecs=30
startretries=30000
autorestart=true
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
process_name=%(program_name)s-%(process_num)s
numprocs=4

[program:task_scheduler]
command=php app/console rabbitmq:consumer -w task_scheduler -l 150 -vv --no-ansi
directory=/www/awardwallet
umask=002
startsecs=30
startretries=30000
autorestart=true
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
process_name=%(program_name)s-%(process_num)s
numprocs=1

[program:scheduler]
command=php app/console aw:scheduler -vv --no-ansi
directory=/www/awardwallet
umask=002
startsecs=30
startretries=50
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/www/awardwallet/app/logs/scheduler.log
stdout_logfile_maxbytes=3000KB
process_name=%(program_name)s-%(process_num)s
numprocs=1
